# scripts/IN4391-server.conf
description "A server node for the IN4391 course"
author "Erwin van Eyk"  

start on runlevel [3]
stop on shutdown  

expect fork

script
    file_name="/opt/IN4391/target/IN4391-server.jar"
    echo "Starting server..."
    if [ -f $file_name ];
    then
      echo "Found jar. Forking process."
      java -jar $file_name > /tmp/IN4391-server.log 2>&1
    else
      echo "Unable to locate jar. Stopped."
    fi
end script
