# scripts/IN4391-server.conf
description "A server node for the IN4391 course"
author "Erwin van Eyk"  

start on runlevel [3]
stop on shutdown  

expect fork  

script
    java -jar /opt/IN4391/target/IN4391-server.jar > /tmp/IN4391-server.log 2>&1
end script
